# dispotrains.webapp/docker/etc/supervisor/conf.d/cron.conf
[program:cron]
command=/usr/sbin/cron -f
directory=/
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/cron.err
stdout_logfile=/var/log/cron.log
